# supervisor.d/tag-generator.conf
[program:tag-generator]
command=/var/www/tag-generator/bin/python init.py
directory=/var/www/tag-generator/application 
environment=
user=nobody
autorestart=true
